# cron-service/cron.d
* * * * * /var/www/cron.sh >> /var/www/cron.log 2>&1
